# apps/maas/s6/root/etc/s6-overlay/s6-rc.d/svc-maas-proxy/run
#!/command/execlineb -P
s6-envuidgid maas
/usr/sbin/squid -d 5 -Y -f /var/lib/maas/maas-proxy.conf
