# apps/maas/root/etc/s6-overlay/s6-rc.d/svc-maas-http/run
#!/command/execlineb -P
until [ -f /var/lib/maas/http/nginx.conf ]; do sleep 5; done
/usr/sbin/nginx -g 'daemon off;' -c /var/lib/maas/http/nginx.conf
