# apps/maas/s6/root/etc/s6-overlay/s6-rc.d/svc-rackd/run
#!/command/execlineb -P
s6-envuidgid maas s6-env prometheus_multiproc_dir=/var/lib/maas/prometheus /usr/sbin/rackd
