# apps/maas/root/etc/s6-overlay/s6-rc.d/init-maas-packages/run
#!/usr/bin/with-contenv bash

if [ -f /usr/bin/apt ]; then ## Ubuntu
  apt -qq update && apt -yq install maas
else
  echo "**** Unsupported OS ****"
  exit 1
fi
