# apps/maas/root/etc/s6-overlay/s6-rc.d/cfg-syslog-wait/run
#!/command/execlineb -P
until [ -f /var/lib/maas/rsyslog.conf ]; do sleep 5; done
