# apps/maas/s6/root/etc/s6-overlay/s6-rc.d/init-regiond/up
if { mkdir -p /etc/maas }
/etc/s6-overlay/s6-rc.d/init-regiond/run
