# apps/universal-flemarr/root/etc/s6-overlay/s6-rc.d/svc-mod-universal-flemarr/run
#!/usr/bin/with-contenv bash

cd /flemarr || exit
exec s6-setuidgid abc python3 run.py
