# apps/maas/root/etc/s6-overlay/s6-rc.d/svc-maas-agent/run
#!/command/execlineb -P
until [ -f /etc/maas/agent.yaml ]; do sleep 5; done
s6-envuidgid -D maas:maas: -- /usr/sbin/maas-agent
