# apps/maas/s6/root/etc/s6-overlay/s6-rc.d/svc-agent/run
#!/command/execlineb -P
s6-envuidgid maas
/usr/sbin/maas-agent
