# apps/maas/s6/root/etc/s6-overlay/s6-rc.d/init-rackd/up
if { mkdir -p /etc/maas }
if { printf 'maas_url:\n- http://localhost:5240/MAAS\n' > /etc/maas/rackd.conf }
chown maas:maas /etc/maas/rackd.conf
