# apps/universal-exportarr/root/etc/s6-overlay/s6-rc.d/svc-mod-universal-exportarr/run
#!/command/execlineb -P
exportarr \
  "$(s6-ps -l -o comm | awk 'NR>1 && !/^(s6-|busybox|bash)/ {print tolower($0); exit}')" \
  --port "${METRICS_PORT:-$((${PORT:-9700}+1))}" \
  --config /config/config.xml \
  --log-level "${LOG_LEVEL:-INFO}" \
  --disable-ssl-verify
