# apps/maas/s6/root/etc/s6-overlay/s6-rc.d/init-maas-proxy/run
#!/command/with-contenv sh
# shellcheck disable=SC1091
source /usr/share/maas/maas-proxy-common.sh
pre_start
