# apps/maas/s6/root/etc/s6-overlay/s6-rc.d/init-maas-proxy/up
if { -f /var/lib/maas/maas-proxy.conf }
/etc/s6-overlay/s6-rc.d/init-maas-proxy/run
