# apps/maas/s6/root/etc/s6-overlay/s6-rc.d/svc-apiserver/run
#!/command/execlineb -P
until [ -f /etc/maas/regiond.conf ]; do sleep 5; done
s6-envuidgid maas /usr/sbin/maas-apiserver
