# apps/maas/s6/root/etc/s6-overlay/s6-rc.d/svc-regiond/run
#!/command/execlineb -P
s6-envuidgid maas s6-env DJANGO_SETTINGS_MODULE=maasserver.djangosettings.settings prometheus_multiproc_dir=/var/lib/maas/prometheus /usr/sbin/regiond
